# docker/temperature/temperature.conf
[supervisord]
nodaemon=true

[program:temperature]
command=python /data/insertTemp.py
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
